# chap2-startup.service
[Unit]
Description=CHAP2 Oracle Cloud Auto-Startup Service
After=network.target docker.service
Wants=network.target docker.service
Requires=docker.service

[Service]
Type=oneshot
User=opc
Group=opc
WorkingDirectory=/home/opc
ExecStart=/home/opc/oracle-cloud-startup.sh
ExecStop=/bin/true
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal
Restart=on-failure
RestartSec=60
TimeoutStartSec=600
TimeoutStopSec=60

# Environment variables
Environment=HOME=/home/opc
Environment=PATH=/usr/local/bin:/usr/bin:/bin

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/opc /var/log /tmp

[Install]
WantedBy=multi-user.target 